# elk/filebeat/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:filebeat]
command = /opt/filebeat -e
priority = 1
autostart = true
autorestart = true
