# caddy/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:caddy]
command = /opt/caddy -conf="/opt/Caddyfile" -agree
priority = 1
autostart = true
autorestart = true
