# oauth-proxy/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:oauth-proxy]
command = /oauth-proxy :8765 /etc/oauth-proxy.toml
priority = 1
autostart = true
autorestart = true
