# couchbase/config/supervisord/supervisord.conf
[supervisord]
nodaemon = true

[program:couchbase]
command = /entrypoint.sh couchbase-server
priority = 1
autostart = true
autorestart = true
