# prom/prometheus/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:prometheus]
command = /run-prometheus.sh
priority = 1
autostart = true
autorestart = true
