# nginx-https-redirector/supervisord.conf
[supervisord]
nodaemon = true
loglevel = debug

[program:nginx]
command = /usr/sbin/nginx
priority = 1
autostart = true
autorestart = true
