# elk/kibana/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:kibana]
command = /opt/kibana/bin/kibana
priority = 1
autostart = true
autorestart = true
