# kong/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:kong]
command = /opt/run-kong.sh
priority = 1
autostart = true
autorestart = true
