# kong/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:kong]
command = /usr/local/bin/kong start
priority = 1
autostart = true
autorestart = true
