# prom/blackbox-exporter/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:blackbox-exporter]
command = /opt/blackbox_exporter -config.file /opt/blackbox_exporter.yml
priority = 1
autostart = true
autorestart = true
