# elk/logstash/supervisord.conf
[supervisord]
nodaemon=true
loglevel = debug

[program:logstash]
command = /opt/logstash/bin/logstash -f /etc/logstash.conf
priority = 1
autostart = true
autorestart = true
